# External/Agent/agent.service
[Unit]
Description=My Home Agent

[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
WorkingDirectory=/home/pi/MyHome/External/Agent
ExecStart=/usr/bin/python3 /home/pi/MyHome/External/Agent/main.py

[Install]
WantedBy=multi-user.target