# External/myhome.service
# copy to /etc/systemd/system/

[Unit]
Description=My Home

[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
WorkingDirectory=/home/pi/MyHome
ExecStart=/home/pi/MyHome/External/Agent/venv/bin/python3 /home/pi/MyHome/Start.py

# This environment variable is necessary when dotnet isn't loaded for the specified user.
# To figure out this value, run 'env | grep DOTNET_ROOT' when dotnet has been loaded into your shell.
Environment=DOTNET_ROOT=/home/pi/.dotnet
# Set default display value (used by vlclib)
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target